# examples/supervisord.conf
[supervisord]
nodaemon=true

[program:example]
command=python /home/spam-ham/decorator_style_task_app.py
user=spam-ham
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true